# ops/wsl/oracle-hermes-vm-dashboard-tunnel-healthcheck.timer
[Unit]
Description=Run Oracle Hermes VM dashboard tunnel health check every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=oracle-hermes-vm-dashboard-tunnel-healthcheck.service

[Install]
WantedBy=timers.target
